# docker/resources/supervisor/conf.d/crond.conf
[program:crond]
command=crond -f -l 8
user=root
autostart=true
autorestart=true