# resources/supervisor/conf.d/proxy-server.conf
[program:proxy-server]
command=/usr/sbin/nginx
priority=1
stopsignal=KILL
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
